# deploy/quadlet/elk-poc-elasticsearch.container.in
[Unit]
Description=ELK POC Elasticsearch
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=elk-poc-elasticsearch
Image=@@ELASTICSEARCH_IMAGE@@
HttpProxy=false
Network=elk-poc.network
PublishPort=@@BIND_ADDRESS@@:@@ELASTICSEARCH_PORT@@:9200
Volume=@@RUNTIME_ROOT@@/config/elasticsearch/elasticsearch.yml:/usr/share/elasticsearch/config/elasticsearch.yml:ro,Z
Volume=@@RUNTIME_ROOT@@/config/elasticsearch/elasticsearch.keystore:/usr/share/elasticsearch/config/elasticsearch.keystore:Z
Volume=@@RUNTIME_ROOT@@/pki/elasticsearch:/usr/share/elasticsearch/config/certs:ro,Z
Volume=@@RUNTIME_ROOT@@/data/elasticsearch:/usr/share/elasticsearch/data:Z
Secret=elk_poc_elastic_password,type=mount,target=elastic_password,uid=1000,gid=0,mode=0400
Environment=ELASTIC_PASSWORD_FILE=/run/secrets/elastic_password
Environment=TZ=UTC
Memory=@@ELASTICSEARCH_MEMORY@@
Ulimit=memlock=-1:-1
Ulimit=nofile=65535:65535
AddCapability=IPC_LOCK
LogDriver=journald
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target
